# dilbert_random.timer
[Unit]
Description=Print random Dilbert's Daily Comic Strip

[Timer]
# Begin of lunch time
OnCalendar=13:00

[Install]
WantedBy=timers.target
