# dilbert_today.timer
[Unit]
Description=Print Dilbert's Daily Comic Strip of today

[Timer]
# Begin of lunch time
OnCalendar=13:00

[Install]
WantedBy=timers.target
